# fotobox.service
[Unit]
Description=Fotobox – Raspberry Pi Photo System
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=%h/Fotobox
ExecStart=%h/Fotobox/start.sh
Restart=on-failure
RestartSec=5
Environment=FOTOBOX_PHOTO_DIR=%h/photos

[Install]
WantedBy=multi-user.target
