# fotobox.service
[Unit]
Description=Fotobox – Raspberry Pi Photo System
After=network.target

[Service]
Type=simple
User=fotobox
WorkingDirectory=/home/fotobox/Fotobox
ExecStart=/home/fotobox/Fotobox/venv/bin/python -m server.app
Restart=always
RestartSec=5
Environment=FOTOBOX_PHOTO_DIR=/home/fotobox/photos

[Install]
WantedBy=multi-user.target
