# fotobox.service
[Unit]
Description=Fotobox – Raspberry Pi Photo System
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/Fotobox
ExecStart=/home/pi/Fotobox/start.sh
Restart=on-failure
RestartSec=5
Environment=FOTOBOX_PHOTO_DIR=/home/pi/photos

[Install]
WantedBy=multi-user.target
